# hypnic.service
[Unit]
Description=Hypnic
After=rescue.target

[Service]
Type=simple
ExecStart=/usr/bin/python /usr/bin/hypnic.py

[Install]
WantedBy=rescue.target
