# sbc/hypnic.service
[Unit]
Description=Hypnic

[Service]
Type=simple
ExecStart=/usr/bin/python /usr/bin/hypnic.py -D

[Install]
WantedBy=basic.target
