# sbc/hypnic.service
[Unit]
Description=Hypnic

[Service]
Type=simple
ExecStart=/usr/bin/hypnic.py -d

[Install]
WantedBy=basic.target